# nomad/supervisor.ini
[program:consul]
command=/usr/local/bin/consul agent -datacenter %(ENV_CONSUL_DATACENTER)s -config-dir /config -client 0.0.0.0 -retry-join %(ENV_CONSUL_SERVER)s -bind 0.0.0.0 -grpc-port=8502 -data-dir /tmp -hcl "enable_central_service_config = true"
autorestart=true
priority=100

[program:docker]
command=/usr/bin/dockerd
autorestart=true
startretries=100
priority=100

[program:nomad]
command=/usr/local/bin/nomad agent -config /etc/nomad.d/config.hcl
autorestart=true
startretries=300
priority=300