# nomad/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
username=username
password=password

[supervisord]
user=root
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
username=username
password=password

[program:consul]
command=/usr/bin/consul agent -config-dir /etc/consul.d/config
autorestart=true
priority=100
stopsignal=SIGTERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:docker]
command=/usr/local/bin/dockerd-entrypoint.sh
autorestart=true
priority=100
#stdout_logfile=/dev/stdout
#stderr_logfile=/dev/stderr
#stdout_logfile_maxbytes=0
#stderr_logfile_maxbytes=0

[program:nomad]
command=/etc/nomad.d/start_nomad.sh
autorestart=true
priority=300
stopsignal=SIGTERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=HTTP_PROXY="",HTTPS_PROXY="",NO_PROXY=""